# crontabs/root
* * * * * /bin/sh /root/periodic/time_log
